# n8n/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root
environment=WEBHOOK_URL="%(ENV_WEBHOOK_URL)s",TIMEZONE="%(ENV_TIMEZONE)s",N8N_PATH="%(ENV_N8N_PATH)s"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
startsecs=1

[program:n8n]
command=bash -c "source /n8n-exports.sh && exec n8n start"
user=node
directory=/home/node
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=5
startsecs=5
stopwaitsecs=30
environment=HOME="/home/node",WEBHOOK_URL="%(ENV_WEBHOOK_URL)s",TIMEZONE="%(ENV_TIMEZONE)s",N8N_PATH="%(ENV_N8N_PATH)s"
